# packaging/systemd/wc-privhelper.service
[Unit]
Description=WiFiCatcher privileged radio helper
Documentation=https://github.com/0xPR3ST1JH0NN7/WiFiCatcher
Requires=wc-privhelper.socket
After=wc-privhelper.socket

[Service]
Type=simple
# Started on demand by wc-privhelper.socket, which passes the listening socket
# as fd 3. install-helper.sh fills in the venv python, install dir and app uid.
ExecStart=/opt/wificatcher/.venv/bin/python -m WiFiCatcher.privileged --idle 300 --peer-uid APP_UID
WorkingDirectory=/opt/wificatcher

# Runs as root: airmon-ng/airodump-ng/aireplay-ng and the NetworkManager restart
# need it. The blast radius is bounded by the tiny, validated request surface,
# and by --peer-uid (only the app user's connections are accepted).
User=root

# Isolate the capture temp files. Stronger sandboxing (ProtectSystem, capability
# bounding, RestrictAddressFamilies) is intentionally left off: the radio tools
# touch the network stack and system services, and the install often lives under
# /home, so those directives (e.g. ProtectHome) break it. Add them only after
# testing on real hardware.
PrivateTmp=true

# No [Install]: this unit is socket-activated, not enabled at boot.
